# letsencrypt-lighttpd.service
[Unit]
Description=Letsencrypt manual renewal

[Service]
Type=oneshot
ExecStart=/usr/local/bin/letsencrypt-lighttpd/renew.sh
